# etc/xdg/autostart/drauger-welcome.desktop
[Desktop Entry]
Name=Drauger Welcome
Version=4.6
Exec=/usr/share/drauger-welcome/tut.py
Comment=Welcome Screen for Drauger OS
Type=Application
Terminal=false
StartupNotify=true
Encoding=UTF-8
Categories=Accessibility;Application;GTK;System;Utility;
